# etc/udev/rules.d/99-usb-mount.rules
ACTION=="add", KERNEL=="sd[b-z][0-9]*", SUBSYSTEM=="block", RUN+="/usr/local/bin/drives_mount.sh %k"